# deploy_tools/gunicorn-systemd.template.service
[Unit]
Description=Gunicorn server for SITENAME

[Service]
Restart=on-failure
User=huy
WorkingDirectory=/home/huy/sites/SITENAME/source
ExecStart=/home/huy/sites/SITENAME/virtualenv/bin/gunicorn \
    --bind unix:/tmp/SITENAME.socket \
    superlists.wsgi:application

[Install]
WantedBy=multi-user.target
